# source/supervisor/dyc_status.conf
[program:dyc_status]
command = /usr/bin/python3 /robot/mq/dyc_status.py
directory = /robot/mq
autostart = yes
autorestart = yes
user = root
stdout_logfile_maxbytes = 200MB
stderr_logfile_maxbytes = 200MB
