# source/supervisor/nmr_status.conf
[program:nmr_status]
command = /usr/bin/python3 /robot/source/robot_observer/nmr_status.py
directory = /robot/source
autostart = yes
autorestart = yes
user = root
stdout_logfile_maxbytes = 200MB
stderr_logfile_maxbytes = 200MB
